# deploy/systemd/ai-daily.service
[Unit]
Description=AI Daily — 自动采集与日报生成
Documentation=https://github.com/chenxianseng001/ai-daily
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=YOUR_USER                     # ← 修改为你的用户名
Group=YOUR_USER
WorkingDirectory=/opt/ai-daily/app
ExecStart=/opt/ai-daily/venv/bin/python3 /opt/ai-daily/app/run_daily.py
StandardOutput=journal
StandardError=journal
LogRateLimitIntervalSec=0
TimeoutStartSec=600

# 从 .env 文件加载环境变量（避免在 systemd 配置中写入明文 Secret）
EnvironmentFile=/opt/ai-daily/.env

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
